# etc/hotplug.d/iface/10-promote_secondaries
/sbin/sysctl net.ipv4.conf.br-lan.promote_secondaries=1