# etc/hotplug.d/iface/30-tunnels
[ "$ACTION" == "ifup" ] && {
        proto=$(uci get network.$INTERFACE.proto)
        [ "$proto" == "gre" ] || [ "$proto" == "grev6" ] && {
                tunipaddr=$(uci get network.$INTERFACE.tunipaddr)
                tunnetmask=$(uci get network.$INTERFACE.tunnetmask)
                [ "$proto" == "grev6" ] && {
                        ip addr add $tunipaddr dev gre6-$INTERFACE
                        ip route add $tunnetmask dev gre6-$INTERFACE
                } || {
                        ip addr add $tunipaddr dev gre4-$INTERFACE
                        ip route add $tunnetmask dev gre4-$INTERFACE
                }
        }
}

