# etc/init.d/td_serverd
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2012 OpenWrt.org

USE_PROCD=1
START=97
SERVICE_DAEMONIZE=1

start_service() {
    for i in $( seq 1 10)
    do
        [ -n "$(ifconfig br-lan  | grep "inet addr")" ] || sleep 1
    done
    stop_service
    procd_open_instance td_serverd
    procd_set_param respawn
    procd_set_param command /bin/td_serverd
    procd_close_instance
}