# usr/lib/opkg/info/igmpproxy.list
/usr/sbin/igmpproxy
/etc/uci-defaults/igmpproxy.defaults
/lib/upgrade/keep.d/igmpproxy
/etc/init.d/iptv
/etc/hotplug.d/iface/50-igmpproxy
/etc/config/iptv
