# usr/lib/opkg/info/libwlan_6x-uci.list
/etc/uci-defaults/wireless
/usr/share/acl.d/dwpald_acl.json
/etc/init.d/wireless
/etc/config/wireless_sdk
