# usr/lib/opkg/info/libwlan_6x-uci.list
/etc/init.d/wireless
/etc/config/wireless_sdk
