# usr/lib/opkg/info/pppdv2.list
/etc/init.d/pppoe-srvv2
/etc/init.d/pppdv2
/usr/bin/pppdv2
/usr/bin/pppoe-srvv2
