# usr/lib/opkg/info/td_ddns.list
/usr/bin/inadyn
/usr/bin/88ip
/usr/bin/phddns
/etc/init.d/ddns
/etc/ddns/ddns.sh
