# usr/lib/opkg/info/td_defen_ctrl.list
/etc/init.d/ddosdefenctl
/usr/bin/ddosdefenctl
