# usr/lib/opkg/info/td_flow_statistic_ctl.list
/etc/init.d/td_flow_statistic_ctl
/bin/td_flow_statistic_ctl
