# usr/lib/opkg/info/td_filter_ctrl.list
/etc/init.d/td_filter_ctrl
/etc/ucloud_reset
/etc/init.d/wlanfilter
/usr/bin/td_filter_ctrl
