# usr/lib/opkg/info/td_netlink_recv_online.list
/usr/bin/td_netlink_recv_online
/etc/init.d/td_netlink_recv_online
