# usr/lib/opkg/info/td_ol_srv.list
/etc/init.d/td_ol_srv
/usr/bin/td_ol_srv
