# usr/lib/opkg/info/td_serverd.list
/etc/sh/icmp_drop_rule.sh
/etc/config/ipv6UserCfg
/bin/td_serverd
/etc/sh/rem_web_rule_on.sh
/etc/hotplug.d/iface/10-promote_secondaries
/etc/hotplug.d/iface/05-check-ipv4-default-route
/etc/init.d/watchdog
/etc/sh/hostapd_acs_schedule
/etc/sh/acs_scan.sh
/etc/sh/iw_scan.sh
/etc/hotplug.d/connectivity/01-ap
/etc/init.d/td_serverd
/usr/lib/libtd_server.so
/etc/hotplug.d/iface/80-ppp-raise-odhcp6c
/etc/hotplug.d/iface/10-updatelan
/etc/config/advance
/etc/sh/iw_custom_post_ssid_cmd.sh
/etc/config/qos
/etc/hotplug.d/iface/60-clean-ipv6-conntrack
/etc/sh/iw_custom_post_radio_cmd.sh
/etc/config/tdsch
/etc/sh/update_firewall_wanip.sh
/etc/sh/restart_wifi
/etc/sh/spi_upgrade.sh
