# usr/lib/opkg/info/td_wan_speed.list
/etc/init.d/td_wan_speed
/bin/td_wan_speed
