# usr/lib/opkg/info/wan_type_probe.list
/usr/sbin/wan_type_probe
/etc/init.d/wan_type_probe
