# usr/lib/opkg/info/wav_api_6x-uci.list
/usr/sbin/band_steering
/etc/init.d/bandsteering
